# scripts/supervisor/notos.conf
[program:notos]
directory = /srv/notos
user = root
command = /srv/notos/scripts/run-gunicorn-django.sh
stdout_logfile = /var/log/supervisor/notos.log
stderr_logfile = /var/log/supervisor/notos.log
